# deploy/roon-controller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.roonlabs.controller</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Library/Application Support/RoonController/dist/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>NODE_ENV</key>
      <string>production</string>
      <key>HOST</key>
      <string>0.0.0.0</string>
      <key>PORT</key>
      <string>3333</string>
      <key>LOG_LEVEL</key>
      <string>info</string>
      <key>ROON_TOKEN_PATH</key>
      <string>/Library/Application Support/RoonController/config/roon-token.json</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Library/Application Support/RoonController</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Library/Logs/RoonController/out.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/RoonController/error.log</string>

    <key>KeepAlive</key>
    <true/>
  </dict>
</plist>
